# seerr/rootfs/etc/services.d/seerr/finish
#!/usr/bin/with-contenv bashio

echo "Seerr stopped"
