# seerr/rootfs/etc/services.d/seerr/run
#!/usr/bin/with-contenv bashio

CONFIG_PATH="/config"

TZ="$(bashio::config 'TZ')"
export TZ

echo "Starting Seerr with TZ=$TZ"

mkdir -p "$CONFIG_PATH"
mkdir -p "$CONFIG_PATH/logs"

# ⚠️ montar /config como /app/config (bind mount real)
if ! mountpoint -q /app/config; then
    rm -rf /app/config
    mkdir -p /app/config
    mount --bind "$CONFIG_PATH" /app/config
fi

exec node dist/index.js
